# etc/service/redis/run
#!/bin/bash

chpst -u redis:redis -o 4192 /usr/bin/redis-server /etc/redis/redis.conf