# etc/service/vhsapi/run
#!/bin/bash

VHSAPI=/opt/api

cd $VHSAPI
exec chpst -u api carton exec $VHSAPI/bin/app.pl
